# gcs/newfiles_skel/etc/profile.d/lxterminal.sh

# LXTerminal needs a SHELL environment variable defined
# and this variable is not defined when a user logs in the system
# by using Active Directory authentication
SHELL="/bin/bash"
